# jupyter/jupyterlab-server/jupyterlab.service
[Unit]
Description=JupyterLab

[Service]
User=alex
Type=simple
ExecStart=/usr/local/bin/start_jupyterlab

[Install]
WantedBy=default.target
